# app_yc/app_yc_file_uploader_worker/etc/service/bi_file_uploader_worker/run
#!/usr/bin/env sh

set -x  # debug.
set -eu

exec \
    file_uploader_worker
